# hosts/wyse/mfsroot/etc/cron.d/empt-mlmmj
* * * * * mlmmj -n -q /usr/local/bin/mlmmj-maintd -F -d /empt/mlmmj_queue
