# hosts/wyse/mfsroot/etc/cron.d/empt-periodic
*	*	*	*	*	root	periodic minutely
8	*	*	*	*	root	periodic hourly
