# empt/common_overlay/etc/cron.d/empt-acme
35 4 * * * root -n -q acme.sh --home /var/db/acme --renew
